# dockerfiles/nodejs/nodejs.ini
[program:nodejs]
directory=/opt/nodejs/
user=root
command=node server.js
stdout_logfile=/var/log/nodejs/supervisor.log
stderr_logfile=/var/log/nodejs/supervisor.error.log
startsecs=1
environment=NODE_ENV=production