# django/dockerfiles/django/django.ini
[program:django]
directory=/opt/django/Django/
environment=PYTHONPATH=/opt/django/Django/
user=root
command=/opt/django/bin/gunicorn mysite.wsgi -w 16 --max-requests 1000 -b 0.0.0.0:8001
stdout_logfile=/var/log/django/supervisor.log
stderr_logfile=/var/log/django/supervisor.error.log
startsecs=1
timeout = 60*60