# django/docker/django/django.ini
[program:django]
directory=/opt/django
user=root
command=gunicorn mysite.wsgi -w 16 --max-requests 1000 -b 0.0.0.0:8001
stdout_logfile=/var/log/django/supervisor.log
stderr_logfile=/var/log/django/supervisor.error.log
startsecs=1
timeout = 60*60